# deploy/gedi-proxy.service
[Unit]
Description=gedi-proxy — GEDI streaming interceptor (EasyWay Testuggine GEDI)
After=network.target
Documentation=https://easyway.it/wiki/guides/gedi-check-cli

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/gedi-check
ExecStart=/usr/local/bin/gedi-proxy --port 8765 --upstream https://api.anthropic.com
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gedi-proxy

# Sicurezza
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
